# templates/launchd-task.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Generic launchd task plist. Generated by /setup → "Setup launchd schedules".

  Placeholders:
    {{LAUNCHD_LABEL}}     — unique label, e.g. com.<user>.<project>-report
    {{PROJECT_PATH}}      — absolute path to project directory
    {{HOME}}              — user's home directory
    {{PROMPT}}            — the Claude prompt to execute (e.g. "/report", "/docs sync --publish")
    {{CALENDAR_INTERVAL}} — XML for one or multiple StartCalendarInterval entries
                            (use array form for weekday-specific schedules)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{LAUNCHD_LABEL}}</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{PROJECT_PATH}}/bin/launchd-runner.sh</string>
        <string>{{PROMPT}}</string>
    </array>

    {{CALENDAR_INTERVAL}}

    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/{{LAUNCHD_LABEL}}.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/{{LAUNCHD_LABEL}}.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>{{HOME}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
